# scripts/initfs/devmgr/rules/default.rules
# WASMOS device-manager bootstrap rules (initfs)
# Loaded from /init/devmgr/rules during early boot.
# udev-style syntax only.
SUBSYSTEM=="pci", ATTR{class}=="0x01", ATTR{subclass}=="0x01", RUN+="system/drivers/ata.wap"
SUBSYSTEM=="block", ATTR{unit}=="0", ENV{MOUNT}="/boot", RUN+="system/drivers/fs_fat.wap"
